# deployment/roles/html_extractor/templates/html_extractor.conf
# This file is maintained by Ansible
# Any change will be lost

[program:html_extractor]
command = java -jar {{ html_extractor_jar_path }}
user = {{ username }}
directory = {{ user_home }}
