# deployment/roles/scraper/templates/scrapyd.conf
# This file is maintained by Ansible
# Any change will be lost

[program:scrapyd]
command={{ python_prefix_path }}/bin/scrapyd
user={{ scrapyd_user }}
directory=/home/scrapyd
